# .config/systemd/user/Pessoal-Mount.service
[Unit]

Description=Montar Google Drive (rclone) - Pasta Pessoal
AssertPathIsDirectory=/home/renan/GDrives/Pessoal
After=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount Pessoal: %h/GDrives/Pessoal \
  --config=/home/renan/.config/rclone/rclone.conf \
  --vfs-cache-mode=full \
  --vfs-cache-max-size=5G \
  --dir-cache-time=12h \
  --poll-interval=1m \
  --allow-other

ExecStop=/bin/fusermount -u /home/renan/GDrives/Pessoal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target

